# mcp-server/kubernetes-mcp-server.service
[Unit]
Description=Kubernetes MCP Server
Documentation=https://github.com/containers/kubernetes-mcp-server
After=network.target
Wants=network.target

[Service]
Type=exec
User=stratus
Group=stratus
WorkingDirectory=/home/stratus/git_projects/k8s-tools/mcp-server
ExecStart=/home/stratus/git_projects/k8s-tools/mcp-server/start-mcp-server.sh npx normal
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kubernetes-mcp-server

# Security settings
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=yes
PrivateDevices=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Environment
Environment="KUBECONFIG=/home/stratus/.kube/config"
Environment="MCP_PORT=3001"
Environment="MCP_LOG_LEVEL=4"

[Install]
WantedBy=multi-user.target
